# conf/systemd/certbot-renewal.service
[Unit]
Description=Certbot Renewal

[Service]
ExecStart=/usr/bin/certbot renew --nginx --post-hook "systemctl reload nginx" 
